# terraform-recipe/payload/available-configs/arm-elf-patch/etc/rc.local
#!/bin/sh
#
# Ensure cowrie restarts on system reboot
#
su -c 'cd ~cowrie/cowrie && bin/cowrie restart' cowrie
start-stop-daemon --user cowrie -b --start --exec /home/ubuntu/scripts/log-processor.sh
