# terraform-recipe/payload/etc/cron.d/cowrie-watchdog
* * * * * cowrie /home/ubuntu/scripts/cowrie-watchdog.sh
