# terraform-recipe/payload/config-files/rc.local
#!/bin/sh
#
# Ensure cowrie restarts on system reboot
#
su -c 'cd ~cowrie/cowrie && bin/cowrie restart' cowrie
